# frame-configurator-frontend.service.example
[Unit]
Description=PiFrame Frontend Web Application
After=network.target

[Service]
Type=simple
User=ben
WorkingDirectory=/home/ben/piframe/piframe-frontend
ExecStart=/usr/bin/npm start
Restart=always
Environment=PORT=8080
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target